# files/daytrader.service
[Unit]
Description=daytrader

[Service]
ExecStart=/opt/daytrader/daytrader-wlp.sh
Type=forking

[Install]
WantedBy=multi-user.target

